# templates/default/mongodb.upstart.erb
# Ubuntu upstart file at /etc/init/mongodb.conf
 
limit nofile 64000 64000

kill timeout 30 # wait 30s between SIGTERM and SIGKILL.

pre-start script
  <%= @dbpath ? "mkdir -p #{@dbpath}" : "" %>

  <% if @raid %>
    vol=$(grep <%= @mnt_point %> /proc/mounts |awk '{print $1}')
    disks=$(mdadm --detail $vol |grep active |awk '{print $7}')

    for n in $vol $disks ; do
      blockdev --setra <%= @setra %> $n
    done
  <% end %>

end script

<% if @type != "shard" %>
start on runlevel [2345]
<% end %>
stop on runlevel [06]

exec sudo -u mongodb <%= @daemon ? "#{@daemon}" : "/usr/bin/mongod" %> --config <%= @configfile ? "#{@configfile}" : "/etc/mongodb.conf" %>
