# templates/default/purge.upstart.erb
# Ubuntu upstart file at /etc/init/oldboy.conf
 
kill timeout 30 # wait 30s between SIGTERM and SIGKILL.

start on runlevel [2345]
stop on runlevel [06]

exec /usr/local/bin/oldboy <%= @when_to_run ? "--when #{@when_to_run}" : "" %> <%= "#{@token}" %>
